# templates/vnc-supervisord.conf.j2
[program:vncserver]
priority=25
directory=/root
command=vncserver :1
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/vnc.log
redirect_stderr=true
stopasgroup=true

